# docs/notes/sw8_cuda_reports/2025-02-09-Jetson-Install/mediamtx.service
[Service]
ExecStart=%h/mediamtx/mediamtx %h/mediamtx/mediamtx.yml
Restart=on-failure
RestartSec=5

[Unit]
Wants=network.target
Description=mediamtx

[Install]
WantedBy=default.target

